# boot/systemd/prometheus-session.target
[Unit]
Description=Prometheus OS Session Target
Documentation=https://prometheus-os.dev
BindsTo=graphical-session.target
